# deploy/systemd/ttyd.service.d/override.conf
[Service]
# Run ttyd as the project user so the launcher script can source the user's .env
# and so an internet-exposed shell is NOT running as root.
User=subhash.thakur.india
Group=subhash.thakur.india
# Fresh Environment — EnvironmentFile= is already set in the base unit.
