# deploy/systemd/factory-paperclip.service
[Unit]
Description=Agentryx Dev Factory - Paperclip (port 3101)
After=network.target postgresql.service

[Service]
Type=simple
User=subhash.thakur.india
Group=subhash.thakur.india
WorkingDirectory=/home/subhash.thakur.india/Projects/paperclip
EnvironmentFile=-/home/subhash.thakur.india/Projects/paperclip/.env
Environment=PATH=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin/node scripts/dev-runner.mjs watch
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
