# deploy/systemd/factory-dashboard.service
[Unit]
Description=Agentryx Dev Factory - Dashboard (Vite, port 5173)
After=network.target

[Service]
Type=simple
User=subhash.thakur.india
Group=subhash.thakur.india
WorkingDirectory=/home/subhash.thakur.india/Projects/pixel-factory-ui
Environment=PATH=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=development
ExecStart=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin/npm run dev -- --host 127.0.0.1 --port 5173
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
